# dvfs/systemd/geopm-apply.path
[Unit]
Description=Watch host GEOPM freq config and apply on change

[Path]
PathChanged=/shared/geopm/freq/%H.conf
PathChanged=/shared/geopm/freq/%H.d
PathExistsGlob=/shared/geopm/freq/%H.d/*.conf
Unit=geopm-apply.service

[Install]
WantedBy=multi-user.target
